# test/cookbooks/awsmetadata/files/default/mockaws.conf
[program:mockawsm]
command=/usr/local/mockaws/bin/python /usr/local/mockaws/mockawsm.py